# scripts/launchd/com.powershop.claude-token-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.powershop.claude-token-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/sync-claude-token.sh</string>
  </array>

  <!-- Run at load (when launchctl bootstraps the agent) and every 2 hours. -->
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>7200</integer>

  <!-- Capture both streams to the same log file so the cause of any error is
       in one place. The launchd agent runs as the user, so the log inherits
       the user's permissions. -->
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/com.powershop.claude-token-sync.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/com.powershop.claude-token-sync.log</string>

  <!-- Inherit PATH so /usr/bin/security and /usr/bin/python3 are found. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
